# infra/production/systemd/zumra-restic-forget.timer
[Unit]
Description=Apply Zümra Restic retention daily

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
