# infra/production/systemd/zumra-backup-full.service
[Unit]
Description=Zümra pgBackRest full backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/srv/zumra/current/infra/production/scripts/physical-backup.sh full
