# infra/production/systemd/zumra-backup-restic.service
[Unit]
Description=Zümra encrypted media and config backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/srv/zumra/current/infra/production/scripts/compose.sh exec -T backup /opt/zumra/scripts/restic-backup.sh
