# infra/production/systemd/zumra-restic-check.timer
[Unit]
Description=Check Zümra Restic repository weekly

[Timer]
OnCalendar=Sun *-*-* 04:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
