# infra/production/systemd/zumra-backup-restic.timer
[Unit]
Description=Run Zümra encrypted backup daily

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
