# infra/production/systemd/zumra-restic-forget.service
[Unit]
Description=Apply Zümra Restic retention policy
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/srv/zumra/current/infra/production/scripts/compose.sh exec -T backup /opt/zumra/scripts/restic-retention.sh forget
