# infra/production/systemd/zumra-backup-wal.timer
[Unit]
Description=Sync Zümra WAL repository every ten minutes

[Timer]
OnCalendar=*:0/10
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
