# infra/production/systemd/zumra-ops-monitor.service
[Unit]
Description=Check Zümra operational health
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/srv/zumra/current/infra/production/scripts/ops-monitor.sh
