# infra/production/systemd/zumra-backup-logical.service
[Unit]
Description=Zümra logical PostgreSQL backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/srv/zumra/current/infra/production/scripts/compose.sh exec -T backup /opt/zumra/scripts/logical-backup.sh
