# infra/production/systemd/zumra-backup-wal.service
[Unit]
Description=Zümra encrypted WAL repository sync
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/srv/zumra/current/infra/production/scripts/compose.sh exec -T backup /opt/zumra/scripts/restic-wal-sync.sh
