# infra/production/systemd/zumra-backup-differential.service
[Unit]
Description=Zümra pgBackRest differential backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/srv/zumra/current/infra/production/scripts/physical-backup.sh diff
